# common/root/etc/profile.d/profile.sh
#!/bin/sh

alias sf="php bin/console"
alias lr="php artisan"